# content/etc/supervisor/conf.d/rundeck.conf
[supervisord]
nodaemon=true

[program:rundeck]
command=/opt/supervisor/rundeck
redirect_stderr=true
